# Docker/data/systemd/cleanup-backups.timer
[Unit]
Description=Run cleanup-backups once a day

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target